# debian/postrm
#!/bin/sh
# postrm script for sftpcloudfs
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postrm> `remove'
#        * <postrm> `purge'
#        * <old-postrm> `upgrade' <new-version>
#        * <new-postrm> `failed-upgrade' <old-version>
#        * <new-postrm> `abort-install'
#        * <new-postrm> `abort-install' <old-version>
#        * <new-postrm> `abort-upgrade' <old-version>
#        * <disappearer's-postrm> `disappear' <overwriter>
#          <overwriter-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    remove)
	echo "Removing sftpcloudfs."
	
	update-rc.d -f sftpcloudfs remove

    ;;
    purge)
	echo "Purging sftpcloudfs."

	if (which deluser && getent passwd sftpcloudfs) > /dev/null 2>&1; then
		deluser --system --quiet sftpcloudfs
	fi	

	if (which delgroup && getent group sftpcloudfs) > /dev/null 2>&1; then
		delgroup --system --quiet sftpcloudfs
	fi

	[ -e /etc/sftpcloudfs.conf ] && rm -f /etc/sftpcloudfs.conf
	[ -e /etc/sftpcloudfs ] && rm -fr /etc/sftpcloudfs
	[ -e /var/run/sftpcloudfs ] && rm -fr /var/run/sftpcloudfs
	[ -e /var/log/sftpcloudfs ] && rm -fr /var/log/sftpcloudfs
	[ -e /etc/logrotate.d/sftpcloudfs ] && rm -f /etc/logrotate.d/sftpcloudfs

    ;;
    upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
    ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
